# systemd/backup-etcd-restic.service
[Unit]
Description=etcd restic encrypt
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/backup/bin/etcd-restic.sh
