# systemd/backup-etcd-restic.timer
[Unit]
Description=6-hourly etcd restic encrypt

[Timer]
OnCalendar=*-*-* 02,08,14,20:00:00
Persistent=true

[Install]
WantedBy=timers.target
