# systemd/backup-master.service
[Unit]
Description=Backup master (NFS+MinIO+metrics+summary email)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/backup/bin/master-backup.sh
