# systemd/backup-integrity.service
[Unit]
Description=Backup integrity check
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/backup/bin/integrity-check.sh
